# images/capi/ansible/roles/providers/files/etc/systemd/system/kubernetes.target
[Unit]
Description=Kubernetes

[Install]
WantedBy=multi-user.target