# images/capi/ansible/roles/providers/files/etc/systemd/system/kube-logrotate.timer
[Unit]
Description=Hourly kube-logrotate invocation

[Timer]
OnCalendar=hourly

[Install]
WantedBy=kubernetes.target