# images/capi/ansible/roles/providers/files/etc/systemd/system/kube-node-configuration.service
[Unit]
Description=Configure kubernetes node
After=kube-node-installation.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/kubernetes/configuration-helper.sh
StandardOutput=journal+console

[Install]
WantedBy=kubernetes.target